# .config/systemd/user/geoclue-agent.service
[Unit]
Description=Geoclue Agent

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent
Restart=on-failure

[Install]
WantedBy=default.target
